# docker/supervisor/supervisord.conf
[program:egg-recipe]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/egg-recipe/artisan queue:work --queue=high,default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/egg-recipe/storage/logs/worker.log
stopwaitsecs=3600